# setup/systemd/breadbrich-auto-deploy.service
[Unit]
Description=Breadbrich Engels auto-deploy (poll origin/main and run safe-deploy when behind)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
# Without HOME, the inline `git config --global ...` calls inside
# safe-deploy.sh fail with "fatal: $HOME not set" and the safe.directory
# workaround can't be written, which trips the ERR trap and triggers a
# rollback (caused PR #36's first auto-deploy to fail). Setting HOME to
# /root puts git's global config in the standard place for root-owned
# operations; per-user (breadbrich) git calls are unaffected because
# `su - breadbrich` sets its own HOME via the login shell.
Environment=HOME=/root
ExecStart=/opt/breadbrich-backups/auto-deploy.sh
StandardOutput=append:/opt/breadbrich/logs/auto-deploy.log
StandardError=append:/opt/breadbrich/logs/auto-deploy.error.log
