# setup/systemd/breadbrich-kb.service
[Unit]
Description=Breadbrich Engels KB dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=breadbrich
Group=breadbrich
WorkingDirectory=/opt/breadbrich
Environment=NODE_ENV=production
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=KB_PORT=8080
Environment=CONTEXT_DIR=/opt/breadbrich/groups/discord_1367176706783445193/context
Environment=USERS_FILE=/opt/breadbrich/kb-ui/users.json
Environment=KB_ADMINS=admin
Environment=KB_SUPERADMINS=admin
Environment=DB_PATH=/opt/breadbrich/store/messages.db
ExecStart=/usr/bin/node /opt/breadbrich/kb-ui/server.mjs
Restart=always
RestartSec=5
StandardOutput=append:/opt/breadbrich/logs/breadbrich-kb.log
StandardError=append:/opt/breadbrich/logs/breadbrich-kb.error.log

[Install]
WantedBy=multi-user.target
