# setup/systemd/breadbrich.service
[Unit]
Description=labor.fun orchestrator
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=breadbrich
Group=breadbrich
WorkingDirectory=/opt/breadbrich
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Non-secret deployment config (NODE_ENV, CREDENTIAL_PROXY_HOST, ...) is in
# /opt/breadbrich/setup/breadbrich-deploy.env (rsynced from setup/ in the
# repo). Secrets remain in /opt/breadbrich/.env (loaded by readEnvFile
# inside Node, not by systemd).
EnvironmentFile=-/opt/breadbrich/setup/breadbrich-deploy.env
ExecStart=/usr/bin/node /opt/breadbrich/dist/index.js
Restart=always
RestartSec=5
StandardOutput=append:/opt/breadbrich/logs/breadbrich.log
StandardError=append:/opt/breadbrich/logs/breadbrich.error.log

[Install]
WantedBy=multi-user.target
