# setup/systemd/breadbrich-kb.service
[Unit]
Description=Breadbrich Engels KB dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=breadbrich
Group=breadbrich
WorkingDirectory=/opt/breadbrich
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# All non-secret deployment config (NODE_ENV, KB_PORT, CONTEXT_DIR,
# USERS_FILE, KB_ADMINS, KB_SUPERADMINS, DB_PATH, ...) is in
# /opt/breadbrich/setup/breadbrich-deploy.env (rsynced from setup/ in the
# repo). The Node process additionally loads secrets from
# /opt/breadbrich/.env via readEnvFile.
EnvironmentFile=-/opt/breadbrich/setup/breadbrich-deploy.env
ExecStart=/usr/bin/node /opt/breadbrich/kb-ui/server.mjs
Restart=always
RestartSec=5
StandardOutput=append:/opt/breadbrich/logs/breadbrich-kb.log
StandardError=append:/opt/breadbrich/logs/breadbrich-kb.error.log

[Install]
WantedBy=multi-user.target
